# s6-rc-poc/service/s6-rc.d/poc/run
#!/command/execlineb -P
cd /
fdmove -c 2 1
/opt/s6-rc-poc
