# umh-core/s6-base/s6-rc.d/umh-core-debug/run
#!/command/execlineb -P
cd /
fdmove -c 2 1
/root/go/bin/dlv --headless --listen=:40000 --api-version=2 --accept-multiclient exec /usr/local/bin/umh-core 