# umh-core/s6-base/s6-rc.d/umh-core/run
#!/command/execlineb -P
cd /
fdmove -c 2 1
exec nice -n -10 /usr/local/bin/umh-core 