# umh-core/s6-base/s6-rc.d/umh-core-log-prepare/up
foreground { mkdir -p /data/logs }
foreground { mkdir -p /data/logs/umh-core }
foreground { mkdir -p /data/redpanda/coredump }
foreground { chmod o+w /data/config.yaml }