# umh-core/s6-base/s6-rc.d/umh-core-log-prepare/up
foreground { mkdir -p /logs }
foreground { mkdir -p /logs/umh-core }
foreground { mkdir -p /data/redpanda/coredump }
foreground { chown -R nobody:nobody /data }
foreground { chown -R nobody:nobody /logs }
foreground { chmod o+w /data/config.yaml }