# umh-core/s6-base/s6-rc.d/umh-core-log-prepare/up
foreground { mkdir -p /data/logs }
foreground { mkdir -p /data/logs/umh-core }
foreground { chown -R nobody:nobody /data }