# umh-core/s6-base/s6-rc.d/redpanda/run
#!/command/execlineb -P
/opt/redpanda/bin/redpanda --redpanda-cfg /opt/redpanda/etc/redpanda/redpanda.yaml