# s6-rc-poc/service/s6-rc.d/poc-log/run
#!/command/execlineb -P
logutil-service /data/logs/s6-rc-poc
