# umh-core/s6-base/s6-rc.d/umh-core-log/run
#!/command/execlineb -P
logutil-service /data/logs/umh-core 